# supervisor.conf
; ═══════════════════════════════════════════════════════════════
; Emare Finance — Supervisor Yapılandırması
; Dosya: /etc/supervisor/conf.d/emare-finance.conf
;
; Tüm arka plan servisleri:
;   1. Queue Worker (2 process) — arka plan işleri
;   2. Laravel Scheduler — zamanlanmış görevler (cron yerine)
;   3. Health Check — watchdog (sunucu kontrolü)
;
; Komutlar:
;   supervisorctl start emare:*
;   supervisorctl stop emare:*
;   supervisorctl restart emare:*
;   supervisorctl status emare:*
; ═══════════════════════════════════════════════════════════════

; ─── QUEUE WORKER ───
; Satış bildirimleri, e-fatura gönderimi, rapor oluşturma vb.
[program:emare-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/emare-finance/current/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --max-jobs=500
directory=/var/www/emare-finance/current
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/emare-finance/shared/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=60
stopsignal=QUIT

; ─── LARAVEL SCHEDULER ───
; Tekrarlayan faturalar, stok uyarıları, rapor zamanlaması
; Crontab'a gerek kalmaz — her dakika schedule:run çalıştırır
[program:emare-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/emare-finance/current/artisan schedule:run --verbose --no-interaction >> /var/www/emare-finance/shared/storage/logs/scheduler.log 2>&1; sleep 60; done"
directory=/var/www/emare-finance/current
autostart=true
autorestart=true
startsecs=0
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/emare-finance/shared/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

; ─── HEALTH CHECK (WATCHDOG) ───
; Her 30 saniyede uygulamayı kontrol eder, düşerse yeniden başlatır
[program:emare-health]
process_name=%(program_name)s
command=/bin/bash /var/www/emare-finance/health-check.sh
directory=/var/www/emare-finance
autostart=true
autorestart=true
startsecs=0
startretries=3
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/emare-finance/shared/storage/logs/health-check.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

; ─── GRUP ───
[group:emare]
programs=emare-queue,emare-scheduler,emare-health
